# deploy/systemd/watchdog.service
[Unit]
Description=Robot Edge Watchdog
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/lib/watchdog
RuntimeDirectory=watchdog
RuntimeDirectoryMode=0755
StateDirectory=watchdog
ExecStart=/usr/local/bin/watchdog -config /etc/watchdog/watchdog.json
Restart=always
RestartSec=2
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
